# RAEConnect.service
[Unit]
Description=RAEConnect Server daemon

[Service]
ExecStart=/RAEConnect/RAEConnectServer

[Install]
WantedBy=multi-user.target